# contrib/github-readme/timeforged-card.timer
[Unit]
Description=TimeForged GitHub Profile Card Sync Timer

[Timer]
OnCalendar=*:0/30
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
